# startup/opcua.conf
[program:daphne]
command=/home/opcua/startup/opcua.sh
user=opcua
stdout_logfile=/home/opcua/logs/opcua.log
redirect_stderr=true
autostart=true
autorestart=true